# backend/supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=/app/scripts/start-fastapi.sh
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr